# ansible/roles/node-elasticsearch/files/lib/systemd/system/caddy.service.d/10-elasticsearch.conf
[Service]
User=elasticsearch
Group=elasticsearch
